# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server
autorestart=true
stdout_logfile=/var/log/redis.log

[program:node]
command=npm start

